# deploy/aws/ec2/gunicorn.service
[Unit]
Description=Gunicorn for Proyecto Lenguaje II
After=network.target

[Service]
User=ec2-user
Group=nginx
WorkingDirectory=/var/www/proyecto
EnvironmentFile=/var/www/proyecto/.env
ExecStart=/var/www/proyecto/.venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/run/gunicorn-proyecto.sock \
    proyecto.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
